# systemd/check-republic.service
[Unit]
Description=Check Republic service

[Service]
User=joe
Type=oneshot
ExecStart=/home/joe/bin/check-republicd-9b6433e
