# systemd/check-republic.timer
[Unit]
Description=Check Republic service timer
StartLimitIntervalSec=0

[Timer]
OnUnitActiveSec=10s
OnBootSec=10s

[Install]
WantedBy=timers.target
